# apps/baram-aer/executor-nitro/scripts/baram-host.service
# Baram TEE Host Service
# Install: sudo cp baram-host.service /etc/systemd/system/
# Enable:  sudo systemctl daemon-reload && sudo systemctl enable baram-host
# Start:   sudo systemctl start baram-host
# Logs:    journalctl -u baram-host -f

[Unit]
Description=Baram TEE Host Server
After=network.target nitro-enclaves-allocator.service
Requires=nitro-enclaves-allocator.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/nasun-monorepo/apps/baram-aer/executor-nitro
Environment=NODE_ENV=production
Environment=USE_VSOCK=true
Environment=ENCLAVE_CID=16
Environment=HOST_PORT=3000
Environment=NVM_DIR=/home/ec2-user/.nvm

# Load .env for settlement variables (contract IDs, private key, etc.)
EnvironmentFile=-/home/ec2-user/nasun-monorepo/apps/baram-aer/executor-nitro/.env

# Use bash to load nvm and run node
ExecStart=/bin/bash -c 'source $NVM_DIR/nvm.sh && exec node dist/host/main.js'

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
